# webssh/Supervisor/daphne.ini
[program:daphne]
directory=/datat/wtxu3/webssh
command=daphne -b 172.31.184.60 -p 8001 --proxy-headers webssh.asgi:application
autostart=true
autorestart=true
stdout_logfile=/tmp/websocket.log
redirect_stderr=true